# root-fs/etc/cron.d/backups

SHELL=/bin/bash
PATH=/root/bin:/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin
# MAILTO=root

backupMysqlArgs="--bz2 --create-dirs --nice --time"

0 0 * * *         root ( backup-mysql.sh $backupMysqlArgs )

0 4,12,20 * * *   root ( \
                          /root/bin/backup-mysql.sh --single $backupMysqlArgs; \
                          /root/bin/backup-rdiff.sh \
                        )
